# examples/si/ansible/roles/ens_workstation/templates/stratio-workstation.service.j2
[Unit]
Description=Stratio Workstation

[Service]
ExecStart=/usr/local/stratio/stratio-workstation

[Install]
WantedBy=default.target